# inittab
# /etc/inittab
init:1:sysinit:/sbin/openrc sysinit
def:1:wait:/sbin/openrc default
# rain:1:respawn:/usr/local/bin/checkra1n -cqv
tty1:1:respawn:/usr/local/bin/checkn1x-surface_welcome
# tty1::respawn:/usr/local/bin/checkn1x_welcome
# tty2::respawn:/usr/local/bin/checkra1n
# tty3::respawn:/usr/local/bin/checkn1x_cliprompt
# tty4::respawn:/usr/sbin/usbmuxd -f
# tty5::respawn:/bin/login -f root
bail:1:ctrlaltdel:/sbin/reboot -f
